# polybot.service
[Unit]
Description=Polymarket Trading Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/polymarket_bot
EnvironmentFile=/root/.env
ExecStart=/usr/bin/python3 /root/polymarket_bot/trading_bot_smart.py
Restart=always
RestartSec=10
StandardOutput=append:/root/polybot/bot.log
StandardError=append:/root/polybot/bot.log

[Install]
WantedBy=multi-user.target
